# debian/install
#!/bin/bash

# HARDN Endpoint Deployment Packaging

set -e

echo "Starting HARDN Endpoint installation..."

# Variables
INSTALL_DIR="/opt/hardn"
SERVICE_FILE="/etc/systemd/system/hardn.service"
REPO_URL="https://github.com/opensource-for-freedom/hardn.git"
LOG_FILE="/var/log/hardn-install.log"

# Logging setup
exec > >(tee -a "$LOG_FILE") 2>&1

# Update and install dependencies
echo "Updating package list and installing dependencies..."
sudo apt-get update
sudo apt-get install -y curl wget git python3 python3-pip

# Create installation directory
echo "Creating necessary directories..."
sudo mkdir -p "$INSTALL_DIR"
sudo chown "$USER":"$USER" "$INSTALL_DIR"

# Clone repository
echo "Cloning HARDN repository..."
if git clone "$REPO_URL" "$INSTALL_DIR"; then
    echo "Repository cloned successfully."
else
    echo "Failed to clone repository. Exiting."
    exit 1
fi

# Install Python dependencies
cd "$INSTALL_DIR"
echo "Installing Python dependencies..."
if pip3 install -r requirements.txt; then
    echo "Python dependencies installed successfully."
else
    echo "Failed to install Python dependencies. Exiting."
    exit 1
fi

# Create systemd service file
echo "Setting up systemd service..."
sudo bash -c "cat > $SERVICE_FILE" <<EOL
[Unit]
Description=HARDN Endpoint Service
After=network.target

[Service]
User=$USER
WorkingDirectory=$INSTALL_DIR
ExecStart=/usr/bin/python3 $INSTALL_DIR/main.py
Restart=always

[Install]
WantedBy=multi-user.target
EOL

# Reload systemd and enable service
echo "Reloading systemd and enabling HARDN service..."
sudo systemctl daemon-reload
sudo systemctl enable hardn.service

# Completion message
echo "Installation complete. You may start HARDN with the following command:"
echo "sudo systemctl start hardn.service"